# systemd/.config/systemd/user/auto-display.service
[Unit]
Description=Auto Display Management
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.config/i3/scripts/auto-display.sh
Environment=DISPLAY=:1
Environment=XAUTHORITY=/run/user/1000/gdm/Xauthority
WorkingDirectory=%h
Nice=19
IOSchedulingClass=3

[Install]
WantedBy=default.target
